# launchd/com.local.mlx.rendered.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>          <string>com.local.mlx</string>
  <key>RunAtLoad</key>      <true/>
  <key>KeepAlive</key>      <true/>
  <key>WorkingDirectory</key><string>/Users/martinfr/LocalClineStack/MacM4LocalAgent</string>
  <key>ProcessType</key>    <string>Interactive</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>   <string>/Users/martinfr/LocalClineStack/MacM4LocalAgent/.venvs/mlx/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>   <string>/Users/martinfr</string>
    <key>USER</key>   <string>martinfr</string>
    <key>TMPDIR</key> <string>/var/folders/6c/pyk36nb12vn1t4ynjrq_hkth0000gn/T/</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/martinfr/LocalClineStack/MacM4LocalAgent/.venvs/mlx/bin/python</string>
    <string>-m</string><string>mlx_lm</string>
    <string>server</string>
    <string>--model</string><string>/Users/martinfr/Documents/GitHub/MacM4LocalAgent/models/mlx-community_Qwen3-Coder-Next-4bit</string>
    <string>--host</string><string>127.0.0.1</string>
    <string>--port</string><string>8081</string>
  </array>
  <key>StandardOutPath</key><string>/Users/martinfr/LocalClineStack/MacM4LocalAgent/.logs/mlx.out.log</string>
  <key>StandardErrorPath</key><string>/Users/martinfr/LocalClineStack/MacM4LocalAgent/.logs/mlx.err.log</string>
</dict>
</plist>
