# launchd/com.local.ollama.rendered.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>          <string>com.local.ollama</string>
  <key>RunAtLoad</key>      <true/>
  <key>KeepAlive</key>      <true/>
  <key>WorkingDirectory</key><string>/Users/martinfr/LocalClineStack/MacM4LocalAgent</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>                   <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>OLLAMA_HOST</key>            <string>127.0.0.1:11434</string>
    <key>OLLAMA_KV_CACHE_TYPE</key>   <string>q4_0</string>
    <key>OLLAMA_FLASH_ATTENTION</key> <string>1</string>
    <key>OLLAMA_NUM_PARALLEL</key>    <string>1</string>
    <key>OLLAMA_KEEP_ALIVE</key>      <string>30m</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/ollama</string>
    <string>serve</string>
  </array>
  <key>StandardOutPath</key><string>/Users/martinfr/LocalClineStack/MacM4LocalAgent/.logs/ollama.out.log</string>
  <key>StandardErrorPath</key><string>/Users/martinfr/LocalClineStack/MacM4LocalAgent/.logs/ollama.err.log</string>
</dict>
</plist>
